# units/docker-network@.service.d/docker-network-driver.conf.j2
# Use the {{ driver }} driver for the network
[Service]
Environment='DOCKER_NETWORK_DRIVER={{ driver }}'
